# docker/supervisord.admin.conf
[supervisord]
nodaemon=true
user=root

[program:web]
directory=/app
command=/bin/sh -c 'exec python3 -m gunicorn -b 0.0.0.0:5000 wsgi:app --workers ${WEB_WORKERS:-1} --threads ${WEB_THREADS:-4} --timeout ${WEB_TIMEOUT:-120} --graceful-timeout ${WEB_GRACEFUL_TIMEOUT:-30} --keep-alive ${WEB_KEEPALIVE:-5} --worker-tmp-dir /dev/shm --error-logfile -'
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
